# mobile/etc/profile.d/elm.sh
export ELM_PROFILE="mobile"
export ELM_MODULES="ctxpopup_copypasteUI>entry/api:datetime_input_popup>datetime/api:access_output_tts>access/api"
