# wearable-circle/HVGA/profile.d/eina.sh
export EINA_LOG_DLOG_ENABLE=1
export EINA_LOG_SYSLOG_ENABLE=0
export EINA_LOG_ABORT=1
export EINA_LOG_ABORT_LEVEL=0
